# systemvm/debian/etc/gro-off.rules
SUBSYSTEM=="net" KERNEL=="eth*" RUN+="/opt/cloud/bin/gro-off.sh %k"
